# serverinfo.service
[Unit]
Description=AutoUpdates Discord

[Service]
ExecStart=nodemon ~/serverinfobot/src/bot.js

[Install]
WantedBy=multi-user.target